# scripts/varnish-healthcheck.service
[Unit]
Description=Varnish healthcheck and switchover

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/varnish-switchover.sh
